# link/.crontab
PATH=/usr/bin:/bin:/usr/local/bin:$HOME/bin

0 11 * * * $HOME/bin/dotfiles p
5 10 * * * $HOME/bin/sep start
0 18 * * * $HOME/bin/sep stop
