# link/.crontab
PATH=/usr/bin:/bin:/usr/local/bin:$HOME/bin

0 11 * * * $HOME/bin/dotfiles p
9 11 * * * sep start
0 18 * * * sep stop
