# link/.crontab
PATH=/usr/bin:/bin:/usr/local/bin

0 11 * * * $HOME/bin/dotfiles p
